# misc/init.d/gentoo/zabbix-agentd
#!/sbin/runscript
#
# Allview agent start/stop script.
#
# Written by A.Tophofen
PATH=/bin:/usr/bin:/sbin:/usr/sbin:/home/allview/bin
DAEMON=/home/allview/bin/allview_agentd
NAME=allview_agentd
DESC="Allview 1.0"
PID=/home/allview/lock/$NAME.pid



opts="${opts} reload"


depend() {
		need net
}

start() {
	ebegin "Starting $DESC: $NAME"
	start-stop-daemon --start --pidfile $PID \
		--exec $DAEMON
	eend
}
stop() {
	ebegin "Stopping $DESC: $NAME"
	start-stop-daemon  --stop  --pidfile $PID \
		--exec $DAEMON
	eend
}
